# mac-mini-config/devstral-2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>devstral-2</string>
    <key>ProgramArguments</key>
      <array>
        <string>/opt/homebrew/bin/llama-server</string>
        <string>--model</string>
        <string>/usr/local/models/mistralai_Devstral-Small-2-24B-Instruct-2512-Q8_0.gguf</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8082</string>
        <string>--n-gpu-layers</string>
        <string>999</string>
        <string>--flash-attn</string>
        <string>on</string>
        <string>--ctx-size</string>
        <string>262144</string>
        <string>--jinja</string>
        <string>--no-prefill-assistant</string>
        <string>--reasoning-format</string>
        <string>none</string>
      </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/llama-server/devstral-2.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/llama-server/devstral-2.log</string>
  </dict>
</plist>
